# packaging/deb/DEBIAN/postinst
#!/bin/sh
set -e

if [ ! -f /etc/kura.env ]; then
  umask 077
  cat > /etc/kura.env <<'EOF'
# KurA environment
# DISCORD_TOKEN=YOUR_TOKEN_HERE
OPUS_CACHE=/var/lib/kura/music_opus
RUST_LOG=warn
EOF
fi

systemctl daemon-reload >/dev/null 2>&1 || true

echo "KurA installed. Edit /etc/kura.env then run: sudo systemctl enable --now kura"
exit 0
